# ansible/roles/adsb/files/adsb.service
[Unit]
Description=ADS-B Chaos Consulting Upload Script

[Service]
ExecStart=/opt/adsb/adsb.sh
Restart=always

[Install]
WantedBy=multi-user.target
