# scripts/adsb.service
[Unit]
Description=adsb

[Service]
ExecStart=/opt/adsb/adsb.sh
Restart=always

[Install]
WantedBy=multi-user.target
