# scripts/ais.service
[Unit]
Description=ais

[Service]
ExecStart=/opt/ais/ais.sh
Restart=always

[Install]
WantedBy=multi-user.target
